# automate/launchd/com.stockdaytrade.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stockdaytrade.scheduler</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/seungjunlee/stock day trade agent/automate/scheduler.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/seungjunlee/stock day trade agent</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/seungjunlee/stock day trade agent/automate/logs/scheduler.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/seungjunlee/stock day trade agent/automate/logs/scheduler.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
